# cfbase/defaults/inittab
# /etc/inittab

# System startup
::sysinit:/sbin/init.d/rcS

# Uncomment the following for login on serial line
#ttyS0::respawn:/sbin/getty 38400 ttyS0

# These are standard terminals on vty[1-4]
tty1::respawn:/sbin/getty 38400 tty1
tty2::respawn:/sbin/getty 38400 tty2
tty3::respawn:/sbin/getty 38400 tty3
tty4::respawn:/sbin/getty 38400 tty4

# Comment above and uncomment below to get
# access to the terminals without a password
#tty1::askfirst:/bin/sh
#tty2::askfirst:/bin/sh
#tty3::askfirst:/bin/sh
#tty4::askfirst:/bin/sh

# System shutdown
::shutdown:/sbin/init.d/shutdown
