# fs_config/crontabs/root
# The system wide crontab is in the root's crontab
# but as this system is not multiuser, it is not a
# big problem (i hope)

# m h dom mon dow	command
25 6	* * *		run-parts /etc/cron.daily
47 6	* * 6		run-parts /etc/cron.weekly
52 6	1 * *		run-parts /etc/cron.monthly
#
